# services/supervisor/supervisor.d/banking.ini
[program:banking]
command=php /mnt/app/CanYinKa-BankIng/artisan start
