# services/supervisor/supervisor.d/sao.ini
[program:soa]
command=php /mnt/app/CanYinKa-Group-Api/artisan start
