# services/supervisor/supervisor.d/user.ini
[program:user]
command=php /mnt/app/CanYinKa-User/artisan start
