# services/supervisor/supervisor.d/user_api.ini
[program:user_api]
command=php /mnt/app/CanYinKa-User-Api/artisan start-Rpc-client 
